# systemd/bess-network.service
[Unit]
Description=BESS Thor sensor network bootstrap
Documentation=file:///home/thor/CLAUDE.md
After=network-online.target NetworkManager.service
Wants=network-online.target
# Rate limit: 3 failures in 24h → unit stays failed. Real safety net against
# boot loops if boot-sensors.sh regresses in the future.
StartLimitIntervalSec=86400
StartLimitBurst=3

[Service]
Type=oneshot
RemainAfterExit=yes
# boot-sensors.sh can take up to ~2 min (60s camera discovery sleep + MikroTik
# probe). Anything beyond 240s and the script is genuinely wedged.
TimeoutStartSec=240
ExecStart=/home/thor/bess/scripts/boot-sensors.sh
# Non-fatal: script can partially fail (e.g. MikroTik off) without failing the service.
SuccessExitStatus=0 1
StandardOutput=journal
StandardError=journal
Restart=no

[Install]
WantedBy=multi-user.target
