# systemd/bess-time-bootstrap.service
[Unit]
Description=BESS Thor cold-boot HTTP time bootstrap
Documentation=file:///home/thor/CLAUDE.md
After=network-online.target
Wants=network-online.target
# Intentionally NO Before= on ptp4l/phc2sys. phc2sys -S 0 will not step PHC1
# even if this unit sets a weird wall clock — worst case slaves lag, they
# don't crash. Keep this unit independent of the PTP chain so a time-fetch
# wedge cannot cascade into sensor net downtime.
StartLimitIntervalSec=86400
StartLimitBurst=3
ConditionPathExists=/usr/local/sbin/bess-time-bootstrap

[Service]
Type=oneshot
RemainAfterExit=yes
TimeoutStartSec=240
ExecStart=/usr/local/sbin/bess-time-bootstrap
SuccessExitStatus=0 1
StandardOutput=journal
StandardError=journal
Restart=no

[Install]
WantedBy=multi-user.target
