# systemd/bess-stack.service
[Unit]
Description=BESS Thor docker compose sensor stack (sequenced launcher)
Documentation=file:///home/thor/CLAUDE.md
# After docker is up and the sensor network is bootstrapped. Wants= (not
# Requires=) so a bess-network failure doesn't block the stack — the launcher
# script uses `up -d --no-deps` and most containers don't need 169.254.100.1
# live to start (the camera SDKs retry discovery internally).
After=docker.service bess-network.service network-online.target
Wants=docker.service bess-network.service network-online.target
# Rate limit: 3 failures in 24h → unit stays failed. Structural safety net
# against Tegra watchdog boot loops.
StartLimitIntervalSec=86400
StartLimitBurst=3

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/thor/bess
User=thor
Group=thor
# The launcher is ~125s of deliberate sleeps + docker compose calls. Allow
# 10 min so a slow image pull or lingering container cleanup doesn't trip
# the timeout on a cold boot.
TimeoutStartSec=600
TimeoutStopSec=180
# 30-second grace before touching docker. Gives SSH + Tailscale enough time
# to come up so a human can intervene from a remote shell if something is
# obviously wrong before the stack starts (touch /tmp/abort-bess-stack and
# reboot inside this window to safely interrupt the boot path).
ExecStartPre=/bin/sleep 30
ExecStart=/home/thor/bess/scripts/bess-stack-up.sh
ExecStop=/usr/bin/docker compose -f /home/thor/bess/docker-compose.thor.yml --profile full down
StandardOutput=journal
StandardError=journal
Restart=no

[Install]
WantedBy=multi-user.target
